# heartbeat.service
[Unit]
Description=run heartbeat

[Service]
Type=oneshot
ExecStart=%h/.local/bin/heartbeat
